# supervisor/messenger-worker.ini
[program:messenger-consume]
command=php /home/shkabo/code/quant-db-seed-be1/bin/console messenger:consume async --time-limit=3600
user=shkabo
numprocs=14
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
